# crontab
# supercronic crontab for meta-agents-v4
# Container TZ is America/Sao_Paulo (set in Dockerfile + fly.toml).
# Format: m h dom mon dow command
#
# Daily traffic campaign for cliente brunobracaioli at 10:00 BRT.
0 10 * * * /app/scripts/run-skill.sh create-traffic-brunobracaioli-campaign

# Daily performance + funnel analysis for cliente brunobracaioli at 08:00 BRT (read-only).
# Covers ALL active campaigns (any objective), full conversion funnel via mcp-meta-ads-b2tech —
# see docs/adr/0025-meta-ads-funnel-analytics-b2tech.md (replaces analytic-traffic-..., ADR 0024).
0 8 * * * /app/scripts/run-skill.sh funnel-analytics-brunobracaioli-campaign

# Daily AI summary for cliente brunobracaioli at 23:30 BRT (consolida o dia p/ o Ultron).
30 23 * * * /app/scripts/run-skill.sh daily-summary-brunobracaioli

# On-demand job queue: every minute, claim+run one pending agent_jobs row enqueued by
# the Ultron dashboard (create/activate campaign). Single-flight via mkdir lock — see
# scripts/poll-agent-jobs.sh and docs/adr/0009-on-demand-agent-jobs-queue.md.
* * * * * /app/scripts/poll-agent-jobs.sh

# Autonomous mode: every minute, claim+tick one DUE autonomous_watches row (Ultron narrates a
# long task's progress by voice while the operator is away). Own mkdir lock, ~90s cadence per
# watch. See scripts/poll-autonomous-watches.sh and docs/adr/0019-ultron-autonomous-mode.md.
* * * * * /app/scripts/poll-autonomous-watches.sh
