# crontab
# supercronic crontab for meta-agents-v4
# Container TZ is America/Sao_Paulo (set in Dockerfile + fly.toml).
# Format: m h dom mon dow command
#
# Daily traffic campaign for cliente brunobracaioli at 10:00 BRT.
0 10 * * * /app/scripts/run-skill.sh create-traffic-brunobracaioli-campaign

# Performance analysis for cliente brunobracaioli, every 3 days at 08:00 BRT (read-only).
# `*/3` em day-of-month = dias 1,4,7,...,28,31 (reinicia no mês — idiom-padrão de "a cada 3 dias").
0 8 */3 * * /app/scripts/run-skill.sh analytic-traffic-brunobracaioli-campaign

# Daily AI summary for cliente brunobracaioli at 23:30 BRT (consolida o dia p/ o Ultron).
30 23 * * * /app/scripts/run-skill.sh daily-summary-brunobracaioli

# On-demand job queue: every minute, claim+run one pending agent_jobs row enqueued by
# the Ultron dashboard (create/activate campaign). Single-flight via mkdir lock — see
# scripts/poll-agent-jobs.sh and docs/adr/0009-on-demand-agent-jobs-queue.md.
* * * * * /app/scripts/poll-agent-jobs.sh
